# backup-smart-onboot.service
[Unit]
Description=Run smart backup shortly after user login

[Service]
Type=oneshot
ExecStart=%h/.local/scripts/backup-smart/backup-smart.sh

