# backup-smart.service
[Unit]
Description=Run backup-smart.sh

[Service]
Type=oneshot
ExecStart=%h/.local/scripts/backup-smart/backup-smart.sh
